# services/seastream-oled.timer
[Unit]
Description=Update Seastream OLED
[Timer]
AccuracySec=1s
OnBootSec=1s
OnUnitActiveSec=1s
[Install]
WantedBy=timers.target