# services/seastream-oled.timer
[Unit]
Description=Update OLED
[Timer]
OnUnitActiveSec=1s
AccuracySec=1s
[Install]
WantedBy=timers.target