# services/seastream-upload-recordings.timer
[Unit]
Description=Analyze and upload Seastream recordings every hour
[Timer]
OnCalendar=00/1:05
[Install]
WantedBy=timers.target